# stage2/06-add-custom-configs/files/wifi-connect.service
[Unit]
Description=Wifi connect script
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/pi/PiBot/Software/pi/wifi.py
Restart=on-failure
RestartSec=15s
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target

